# airthingswave-mqtt.timer
#
# Scheduler for airthingswave-mqtt service
#

[Unit]
Description=Send data from Airthings Waves to MQTT regularly

[Timer]
OnCalendar=hourly
